# scripts/health-check.cron
# PushPulse emulator health check — every 5 minutes
*/5 * * * * /home/pushpulse/PushPulse/scripts/health-check.sh
